# deploy/gmag-webui-update.service.in
[Unit]
Description=Check origin and auto-update gmag_webui, restart dashboard on change
Documentation=https://github.com/HamSCI/gmag_webui
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=GMAG_REPO=@REPO@
Environment=GMAG_SERVICE=gmag-webui.service
ExecStart=/usr/local/bin/gmag-webui-update.sh
