# deploy/ws-mqtt-bridge.service.in
[Unit]
Description=HamSCI mag-usb WebSocket to MQTT bridge
After=network-online.target mosquitto.service mag-usb.service
Wants=network-online.target
Requires=mosquitto.service

[Service]
Type=simple
User=@USER@
WorkingDirectory=@BRIDGE_DIR@
EnvironmentFile=@BRIDGE_DIR@/bridge.env
ExecStart=/usr/bin/python3 @BRIDGE_DIR@/mag_mqtt.py
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
